# utils/cron/log.cron
* */23 * * * /home/web/blog/shell/freshlog.sh